# infra/systemd/hermes-gateway.service
[Unit]
Description=Hermes Agent gateway (OEM alert -> LLM RCA -> Google Chat)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
# Adjust user/paths to your VPS layout (placeholders).
User=hermes
Group=hermes
WorkingDirectory=/opt/hermes
EnvironmentFile=/opt/hermes/.env.runtime
ExecStart=/opt/hermes/.venv/bin/uvicorn hermes.main:app --host 0.0.0.0 --port 2020
Restart=on-failure
RestartSec=5
# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
