# central-vm/loki/loki.service
[Unit]
Description=Grafana Loki (log store — Luồng B)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=loki
Group=loki
ExecStart=/usr/local/bin/loki -config.file=/etc/loki/loki-config.yaml
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=full
ReadWritePaths=/var/lib/loki

[Install]
WantedBy=multi-user.target
